# 06-NimbleX/etc/rc.local
#!/bin/bash
# Here we add whatever we want to execute at boot time
