# 06-NimbleX/lib/systemd/system/ntp.service
[Unit]
Description=NTP

[Service]
ExecStart=/usr/sbin/ntpd -n

[Install]
WantedBy=multi-user.target
