# 06-NimbleX/etc/profile.d/beautify.sh
#!/bin/bash
alias grep="grep --color=auto"
alias df="df -x squashfs"
