# 06-NimbleX/usr/lib/systemd/system/avahi-daemon.socket
[Unit]
Description=Avahi mDNS/DNS-SD Stack Activation Socket

[Socket]
ListenStream=/run/avahi-daemon/socket

[Install]
WantedBy=sockets.target
