# 06-NimbleX/usr/lib/systemd/system/nimblex-refresh-caches.service
[Unit]
Description=NimbleX: validate merged-FS caches (ldconfig, GSettings, MIME, ...)
# Bundles ship pre-merged caches built in their AUFS chroot, so this is
# normally a few stat() calls and exits in <50ms. NimbleX autologins
# straight into GUI, so we block the graphical stack to be safe in the
# rare case a rebuild is actually needed.
DefaultDependencies=no
After=local-fs.target systemd-tmpfiles-setup.service
Before=basic.target display-manager.service startx.service graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nimblex-refresh-caches
TimeoutStartSec=120

[Install]
WantedBy=basic.target
