# 06-NimbleX/usr/lib/systemd/system/nimblex-refresh-caches.service
[Unit]
Description=NimbleX: validate merged-FS caches (ldconfig, GSettings, MIME, ...)
Documentation=man:glib-compile-schemas(1) man:update-mime-database(1)
# Bundles ship pre-merged caches built in their AUFS chroot, so this is
# normally a few stat() calls and exits in <50ms. NimbleX autologins
# straight into GUI, so we block the graphical stack to be safe in the
# rare case a rebuild is actually needed.
Before=display-manager.service startx.service graphical.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/nimblex-refresh-caches
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
