# 06-NimbleX/usr/lib/systemd/system/virtqemud-admin.socket
[Unit]
Description=Libvirt QEMU daemon admin socket
Documentation=man:virtqemud(8)

[Socket]
ListenStream=/var/run/libvirt/virtqemud-admin-sock
Service=virtqemud.service

[Install]
WantedBy=sockets.target
