# 06-NimbleX/usr/lib/systemd/system/zram-swap.service
[Unit]
Description=Compressed RAM swap (zram)
Documentation=man:zram(4) https://docs.kernel.org/admin-guide/blockdev/zram.html
DefaultDependencies=no
After=systemd-modules-load.service
Before=swap.target shutdown.target
Conflicts=shutdown.target
# Don't even try where zram is impossible or unwanted.
ConditionPathExists=/sys/class/zram-control
ConditionVirtualization=!container

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/zram-swap
ExecStart=/usr/libexec/zram-swap start
ExecStop=/usr/libexec/zram-swap stop

# Minor hardening — no new privs, no namespace-creating options: a
# private mount namespace makes the kernel record the swap device under
# a path that vanishes when the oneshot exits, breaking `swapoff -a`.
# We intentionally keep CAP_SYS_ADMIN because swapon(2) requires it.
NoNewPrivileges=yes
ProtectClock=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
